# fixtures/packaging/debian-lxd-snap/etc/init/lxe.conf
description "LXE is a shim of the Kubernetes Container Runtime Interface for LXD"

start on runlevel [2345]
stop on shutdown

#respawn
#umask 022
#chroot /
#chdir /
#limit msgqueue <softlimit> <hardlimit>
#limit nice <softlimit> <hardlimit>
#limit rtprio <softlimit> <hardlimit>
#limit sigpending <softlimit> <hardlimit>
#setuid lxe
#setgid lxe

script
    source /etc/default/lxe
    echo $$ > /var/run/lxe.pid
    exec /usr/bin/lxe --logfile /var/log/lxe.log $OPTIONS
end script
